# apps/widerspruch-jobcenter/scripts/com.boardbot.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.boardbot.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alexanderdeibel/fintutto-ecosystem/apps/widerspruch-jobcenter/scripts/start-bot.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/alexanderdeibel/fintutto-ecosystem/apps/widerspruch-jobcenter/logs/bot-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alexanderdeibel/fintutto-ecosystem/apps/widerspruch-jobcenter/logs/bot-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/alexanderdeibel/fintutto-ecosystem/apps/widerspruch-jobcenter</string>
</dict>
</plist>
